# deploy/systemd/news-api.service
[Unit]
Description=Niuma News API + Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/news-aggregator/server
EnvironmentFile=/opt/news-aggregator/.env
ExecStart=/usr/bin/node --import tsx src/index.ts
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
